# src/main/rpm/scriptlets/post-install.sh
#!/bin/sh

/sbin/chkconfig --add wm-config-server
